# system/profile.d/pkg-update.sh
pacman -Syu --noconfirm

ORPHANS=$(pacman -Qdtq)
pacman -Rns "$ORPHANS" --noconfirm
